# files/etc/service/mongo-log/run
#!/bin/bash
tail --retry --follow /config/mongodb/log/mongodb.log
