# files/etc/service/mongo/run
#!/bin/bash
exec setuser system mongod --config /etc/mongodb.conf
